# scripts/com.chekgonbin.dbd-monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyLists-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chekgonbin.dbd-monthly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/visarutsankham/chekgonbin-mvp-1/scripts/run-dbd-monthly.sh</string>
    </array>

    <!-- Run on the 1st of every month at 03:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- If the Mac was asleep at the scheduled time, run it on next wake -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/visarutsankham/chekgonbin-mvp-1/scripts/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/visarutsankham/chekgonbin-mvp-1/scripts/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/visarutsankham/chekgonbin-mvp-1</string>
</dict>
</plist>
